# dockerfiles/config/supervisor/queue.ini
[program:queue]
process_name=%(program_name)s_%(process_num)02d
command=php /home/app/sms/artisan horizon
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/home/app/sms/storage/logs/queue.log
